# luci-app-webdroidx/etc/init.d/webdroidx
#!/bin/sh /etc/rc.common

START=99
EXTRA_COMMANDS="status"

BIN="/usr/share/webdroidx/webdroidx"
WORKDIR="/usr/share/webdroidx"
LOGDIR="/var/log/webdroidx"
LOGFILE="webdroidx.log"

start() {
    screen -S webdroidx -X quit 2>/dev/null
    
    server_port=$(uci -q get webdroidx.config.server_port)
    video_bit_rate=$(uci -q get webdroidx.config.video_bit_rate)
    device=$(uci -q get webdroidx.config.device)
    
    [ -z "$server_port" ] && server_port="8000"
    [ -z "$video_bit_rate" ] && video_bit_rate="1024000"
    
    mkdir -p $LOGDIR
    > $LOGDIR/$LOGFILE
    
    if [ ! -f "$BIN" ]; then
        echo "$(date): ERROR - WebDroidX binary not found at $BIN" >> $LOGDIR/$LOGFILE
        return 1
    fi
    
    cd $WORKDIR
    
    CMD="$BIN -p $server_port --video_bit_rate $video_bit_rate"
    if [ -n "$device" ] && [ "$device" != "" ]; then
        CMD="$CMD --device $device"
    fi
    
    screen -dmS webdroidx bash -c "$CMD >> $LOGDIR/$LOGFILE 2>&1"
    
    sleep 2
    if screen -list | grep -q "webdroidx"; then
        echo "$(date): WebDroidX started successfully on port $server_port" >> $LOGDIR/$LOGFILE
        return 0
    else
        echo "$(date): ERROR - WebDroidX failed to start" >> $LOGDIR/$LOGFILE
        return 1
    fi
}

stop() {
    if screen -list | grep -q "webdroidx"; then
        screen -S webdroidx -X quit 2>/dev/null
        sleep 1
        echo "$(date): WebDroidX stopped" >> $LOGDIR/$LOGFILE
        return 0
    else
        echo "$(date): WebDroidX was not running" >> $LOGDIR/$LOGFILE
        return 1
    fi
}
restart() {
    echo "$(date): WebDroidX restarting..." >> $LOGDIR/$LOGFILE
    stop
    sleep 2
    start
}

status() {
    if screen -list | grep -q "webdroidx"; then
        echo "WebDroidX is running"
        return 0
    else
        echo "WebDroidX is not running"
        return 1
    fi
}
